# src/infra/packaging/deb/systemd/ai-distro-agent.socket
[Unit]
Description=AI Distro Agent IPC Socket

[Socket]
ListenStream=/run/ai-distro/agent.sock
SocketMode=0660
DirectoryMode=0755
SocketUser=root
SocketGroup=ai-distro

[Install]
WantedBy=sockets.target
