# src/infra/packaging/deb/systemd/ai-distro-core.service
[Unit]
Description=AI Distro Core
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/ai-distro-core
Restart=on-failure
RestartSec=2
RuntimeDirectory=ai-distro
RuntimeDirectoryMode=0755
StateDirectory=ai-distro-core
LogsDirectory=ai-distro-core
User=ai-distro
Group=ai-distro
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
