# src/infra/packaging/deb/systemd/ai-distro-voice.service
[Unit]
Description=AI Distro Voice Pipeline
After=sound.target

[Service]
Type=simple
ExecStart=/usr/bin/ai-distro-voice
Restart=on-failure
RestartSec=2
RuntimeDirectory=ai-distro
RuntimeDirectoryMode=0755
StateDirectory=ai-distro-voice
LogsDirectory=ai-distro-voice
User=ai-distro
Group=ai-distro
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
